# mac/com.claude.pager.pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job for claude-pull. Runs every 60 seconds, syncing any
  newly-completed agent sessions into ~/ClaudeRuns/.

  Install with mac/install_launchd.sh — that script substitutes the
  installed path of claude-pull into a copy at:
    ~/Library/LaunchAgents/com.claude.pager.pull.plist
  then `launchctl bootstrap`s it.

  The placeholder __CLAUDE_PULL_PATH__ MUST be replaced before
  loading; otherwise launchd rejects the file silently.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude.pager.pull</string>

  <key>ProgramArguments</key>
  <array>
    <string>__CLAUDE_PULL_PATH__</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/claude-pull.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/claude-pull.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>

  <key>ProcessType</key>
  <string>Background</string>

  <key>ThrottleInterval</key>
  <integer>30</integer>
</dict>
</plist>
